# image/centos_node:lastet/supervisord.d/init-app.camtech.com.co.conf
[program:react-runServer]
command= npm start

directory=/var/www/html
stopasgroup=true
numprocs=1
autostart=true
autorestart=true
startretries=1

